# doc/disablegpe08/disablegpe08.service
[Unit]
Description=Disable GPE 08

[Service]
Type=simple
ExecStart=/usr/local/bin/disablegpe08
